# monero/.footprint
drwxr-xr-x	root/root	etc/
-rwxr-xr-x	root/root	etc/monerod.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/monerod
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/monero-blockchain-blackball
-rwxr-xr-x	root/root	usr/bin/monero-blockchain-export
-rwxr-xr-x	root/root	usr/bin/monero-blockchain-import
-rwxr-xr-x	root/root	usr/bin/monero-blockchain-usage
-rwxr-xr-x	root/root	usr/bin/monero-gen-trusted-multisig
-rwxr-xr-x	root/root	usr/bin/monero-wallet-cli
-rwxr-xr-x	root/root	usr/bin/monero-wallet-rpc
-rwxr-xr-x	root/root	usr/bin/monerod
